# doc/proxy/crontab
PATH=/home/ipin/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

@reboot mkdaemon -r crazyproxy -Cb www.55boom.net:3400 :8800 >> ~/log 2>&1
